# etc/profile.d/mwzoo.sh
export MWZOO_HOME=/opt/mwzoo
export PATH="$PATH:$MWZOO_HOME"
